# scripts/fool-server-restore.service
[Unit]
Description=Fool Server Restore
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/fool-server-restore
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=4000
Environment=SERVE_WEB=1
ExecStart=/usr/bin/npm run start:prod
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
